# app/alpine-php/74/supervisord.conf
[supervisord]
user=root
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:php-fpm]
# 这里就是你提到的命令！
# -F 表示在前台运行，这样 Supervisor 才能监控它的状态
command=php-fpm -F
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=5

# 如果你有异步队列（比如 Laravel Queue），可以紧接着写在这里
# [program:php-worker]
# command=php /var/www/html/artisan queue:work
# autostart=true
# autorestart=true
# user=www-data
